# deploy/supervisord.conf
; supervisord.conf — ZidaBot process supervisor (alternative to systemd)
; Usage (macOS / Docker / non-systemd Linux):
;   pip install supervisor
;   supervisord -c deploy/supervisord.conf
;   supervisorctl -c deploy/supervisord.conf status

[supervisord]
nodaemon=false
logfile=/var/log/zidabot/supervisord.log
pidfile=/var/run/zidabot/supervisord.pid
childlogdir=/var/log/zidabot

[unix_http_server]
file=/var/run/zidabot/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/zidabot/supervisor.sock

; ---- Brain ----
[program:brain]
command=/usr/local/bin/brain --config /etc/zidabot/brain.conf
autostart=true
autorestart=true
startsecs=3
startretries=10
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/var/log/zidabot/brain.log
stderr_logfile=/var/log/zidabot/brain.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5

; ---- PoP processes (one per venue) ----
[program:pop_binance]
command=/usr/local/bin/pop --config /etc/zidabot/pop_binance.conf
autostart=true
autorestart=true
startsecs=3
startretries=10
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/var/log/zidabot/pop_binance.log
stderr_logfile=/var/log/zidabot/pop_binance.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3

[program:pop_okx]
command=/usr/local/bin/pop --config /etc/zidabot/pop_okx.conf
autostart=true
autorestart=true
startsecs=3
startretries=10
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/var/log/zidabot/pop_okx.log
stderr_logfile=/var/log/zidabot/pop_okx.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3

[program:pop_bybit]
command=/usr/local/bin/pop --config /etc/zidabot/pop_bybit.conf
autostart=true
autorestart=true
startsecs=3
startretries=10
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/var/log/zidabot/pop_bybit.log
stderr_logfile=/var/log/zidabot/pop_bybit.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3

[program:pop_bitget]
command=/usr/local/bin/pop --config /etc/zidabot/pop_bitget.conf
autostart=true
autorestart=true
startsecs=3
startretries=10
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/var/log/zidabot/pop_bitget.log
stderr_logfile=/var/log/zidabot/pop_bitget.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3

[program:pop_kucoin]
command=/usr/local/bin/pop --config /etc/zidabot/pop_kucoin.conf
autostart=true
autorestart=true
startsecs=3
startretries=10
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/var/log/zidabot/pop_kucoin.log
stderr_logfile=/var/log/zidabot/pop_kucoin.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3

; ---- Process group (control all PoPs together) ----
[group:pops]
programs=pop_binance,pop_okx,pop_bybit,pop_bitget,pop_kucoin
